# rel/ems-bus_amd64/etc/systemd/system/ems-bus.service
[Unit]
Description=ems_bus
After=network.target

[Service]
Type=forking
User=agilar
Group=users
ExecStart=ems-bus start

[Install]
WantedBy=multi-user.target
