# rel/deb/ems-bus-Ubuntu-yakkety_amd64/usr/lib/ems-bus/lib/ems_bus-1.0.7/priv/systemd/ems_bus.service
[Unit]
Description=ems_bus
After=network.target

[Service]
Type=forking
User=agilar
Group=users
ExecStart=/var/opt/ems_bus/bin/ems_bus start

[Install]
WantedBy=multi-user.target
